# unstrain.service
[Unit]
Description=20-20-20

[Service]
ExecStart=/usr/local/bin/unstrain
Restart=always

[Install]
WantedBy=default.target
